# systemd-units/echo.socket
[Unit]
Description=echo socket

[Socket]
ListenStream=22222
Accept=true