# systemd-units/loggertest.service
[Unit]
Description=Example Test service

[Service]
Type=oneshot
ExecStart=/usr/bin/logger -p local3.debug I\'m a logger