# systemd-units/test2.target
[Unit]
Description=test 2
Wants=test1.target